# home/.config/autostart/kmonad.desktop
[Desktop Entry]
Name=Restart KMonad so it loads user PATH
Type=Application
Exec=systemctl --user restart kmonad
Terminal=false

