# home/.config/autostart/kmonad.desktop
[Desktop Entry]
Name=Load xkbcomp keyboard map modifications
Type=Application
Exec=systemctl --user restart kmonad
Terminal=false

