# home/.config/systemd/user/xkbcomp.service
[Unit]
Description=Load XKB modifications

[Service]
Type=simple
ExecStart=%h/bin/xkbcomp-load-mods

[Install]
WantedBy=default.target
